# container/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

[program:ws_server]
command=python3 /container/ws_server.py
autostart=true
autorestart=true
startretries=5
stdout_logfile=/tmp/ws_server.log
stderr_logfile=/tmp/ws_server.log
priority=1

[program:file_watcher]
command=python3 /container/file_watcher.py
autostart=true
autorestart=true
startretries=3
stdout_logfile=/tmp/file_watcher.log
stderr_logfile=/tmp/file_watcher.log
priority=2

[program:heartbeat]
command=python3 /container/heartbeat.py
autostart=true
autorestart=false
stdout_logfile=/tmp/heartbeat.log
stderr_logfile=/tmp/heartbeat.log
priority=3
