# container/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

[program:ws_server]
command=python3 /container/ws_server.py
autostart=true
autorestart=true
startretries=5
stdout_logfile=/tmp/ws_server.log
stderr_logfile=/tmp/ws_server.log
priority=1

[program:file_watcher]
command=python3 /container/file_watcher.py
autostart=true
autorestart=true
startretries=3
stdout_logfile=/tmp/file_watcher.log
stderr_logfile=/tmp/file_watcher.log
priority=2

[program:heartbeat]
command=python3 /container/heartbeat.py
autostart=true
autorestart=false
stdout_logfile=/tmp/heartbeat.log
stderr_logfile=/tmp/heartbeat.log
priority=3
; when heartbeat exits (session timeout), stop everything
stopasgroup=true
killasgroup=true

[eventlistener:stop_on_heartbeat_exit]
command=python3 -c "
import sys, os
from supervisor.childutils import listener
while True:
    headers, payload = listener.wait(sys.stdin, sys.stdout)
    if 'PROCESS_STATE_EXITED' in headers.get('eventname',''):
        data = dict(x.split(':') for x in payload.split())
        if data.get('processname') == 'heartbeat' and data.get('expected') == '1':
            os.system('kill -TERM \$(cat /tmp/supervisord.pid)')
    listener.ok(sys.stdout)
"
events=PROCESS_STATE_EXITED
